# server_deployment/fts_server_watcher.path
[Path]
PathModified=/srv/feel-the-streets/fts_server

[Install]
WantedBy=multi-user.target