# server_deployment/fts_worker_watcher.path
[Path]
PathModified=/srv/feel-the-streets/worker
PathModified=/srv/feel-the-streets/*.yml

[Install]
WantedBy=multi-user.target